# utils/get-vnc-users/DEBIAN/postinst
#!/bin/bash
systemctl enable get-vnc-users
systemctl restart get-vnc-users