# utils/get-vnc-users/DEBIAN/prerm
#!/bin/bash
systemctl stop get-vnc-users
systemctl disable get-vnc-users