# utils/get-current-ip/DEBIAN/postinst
#!/bin/bash
systemctl enable get-current-ip
systemctl restart get-current-ip