# utils/get-current-ip/DEBIAN/prerm
#!/bin/bash
systemctl stop get-current-ip
systemctl disable get-current-ip