# utils/save-photo/DEBIAN/postinst
#!/bin/bash
systemctl enable save-photo
systemctl restart save-photo