# utils/save-photo/DEBIAN/prerm
#!/bin/bash
systemctl stop save-photo
systemctl disable save-photo